# share/89-usbcopy.rules
KERNEL=="sd[a-z]?", SUBSYSTEM=="block", ACTION=="add", RUN+="/usr/local/bin/import_music.sh %k"
